# conf_BBB/autostart.service
[Unit]
Description=Copy .profile

[Service]
ExecStart=/etc/init.d/autostart start
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
